# deploy/systemd/hoberadius-license-panel-whatsapp-drain.service.example
[Unit]
Description=HobeRadius License Panel - WhatsApp queue drain (oneshot)
After=network.target

[Service]
Type=oneshot
User=licensepanel
Group=licensepanel
WorkingDirectory=/opt/hoberadius-license-panel
EnvironmentFile=/etc/hoberadius-license-panel/license-panel.env
ExecStart=/opt/hoberadius-license-panel/.venv/bin/flask --app wsgi whatsapp-drain
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/hoberadius-license-panel/instance /var/backups/hoberadius-license-panel

[Install]
WantedBy=multi-user.target
